# 50-mouse.rules
ACTION=="bind", SUBSYSTEMS=="misc", RUN+="/home/xu/workspace/debian12-i3wm/fixmouse.sh"
